# roles/healthCheck/templates/healthCheck.timer.j2
[Unit]
Description=Run Health Check Script Every 1 Second

[Timer]
OnUnitActiveSec=1s
Unit=healthCheck.service

[Install]
WantedBy=timers.target
