# roles/healthCheck/templates/healthCheck.service.j2
[Unit]
Description=Health Check Script
After=network.target

[Service]
ExecStart=/usr/local/bin/healthCheck.sh
StandardOutput=journal
StandardError=journal
