# roles/healthCheck/templates/webhookHandler.timer.j2
[Unit]
Description=Run Webhook Handler Script Every 1 Second

[Timer]
OnUnitActiveSec=1s
Unit=webhookHandler.service

[Install]
WantedBy=timers.target